# files/system/etc/systemd/system/system-flatpak-setup.service.d/10-firstboot-only.conf
[Unit]
# Run only once per installation; marker file is created in ExecStartPost.
ConditionPathExists=!/var/lib/cinnamon-flatpak-setup/done

[Service]
# Create state directory on demand and mark completion after successful run.
StateDirectory=cinnamon-flatpak-setup
ExecStartPost=/usr/bin/touch /var/lib/cinnamon-flatpak-setup/done
