# files/system/etc/systemd/system/bluefin-brew-sync.service
[Unit]
Description=Install Bluefin curated Homebrew packages
# Require network for fallback container pull path.
After=network-online.target
Wants=network-online.target
# Skip if sync has already completed.
ConditionPathExists=!/var/lib/bluefin-brew-sync/done

[Service]
Type=oneshot
# Default desktop user UID; script accepts UID or username.
Environment=TARGET_USER=1000
ExecStart=/usr/bin/bash /usr/libexec/cinnamon-ublue/bluefin-brew-sync.sh
# Persist marker directory used by ConditionPathExists/state file.
StateDirectory=bluefin-brew-sync
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
